# pi/systemd/robo-arm.service
[Unit]
Description=Robo Arm Pi Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_ROOT__/pi
Environment=PYTHONUNBUFFERED=1
ExecStart=__REPO_ROOT__/pi/scripts/run_pi_stack.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
